# nodes/broker/mbus-radio/templates/mbus_radio.service.in
[Unit]
Description=mBus (radio) server v${${PROJECT_NAME}_VERSION_MAJOR}.${${PROJECT_NAME}_VERSION_MINOR}.${${PROJECT_NAME}_VERSION_PATCH} (project ${PROJECT_NAME})
After=network.target

[Service]
Type=simple
ExecStart=/opt/smf/mbus_radio -C/opt/smf/mbus_radio_v${${PROJECT_NAME}_VERSION_MAJOR}.${${PROJECT_NAME}_VERSION_MINOR}.cfg
Restart=always
RestartSec=12

[Install]
WantedBy=multi-user.target
