# nodes/http/templates/http.service.in
[Unit]
Description=http server v${${PROJECT_NAME}_VERSION_MAJOR}.${${PROJECT_NAME}_VERSION_MINOR}.${${PROJECT_NAME}_PATCH_LEVEL} (project ${PROJECT_NAME})
After=network.target

[Service]
Type=simple
ExecStart=/opt/node/http -C/opt/node/http_v${${PROJECT_NAME}_VERSION_MAJOR}.${${PROJECT_NAME}_VERSION_MINOR}.cfg
Restart=always
RestartSec=12

[Install]
WantedBy=multi-user.target
