# nodes/ipt/segw/templates/segw.service.in
[Unit]
Description=smart energy gateway v${${PROJECT_NAME}_VERSION_MAJOR}.${${PROJECT_NAME}_VERSION_MINOR}.${${PROJECT_NAME}_VERSION_PATCH} (project ${PROJECT_NAME})
#
# use this dependency when running on the same machine
#
#Wants=node-ipt.service
#After=network.target node-ipt.service
After=network.target

[Service]
Type=simple
ExecStart=/opt/smf/segw -C/opt/smf/segw_v${${PROJECT_NAME}_VERSION_MAJOR}.${${PROJECT_NAME}_VERSION_MINOR}.cfg
Restart=always
RestartSec=12

[Install]
WantedBy=multi-user.target
