# nodes/ipt/gateway/templates/gateway.service.in
[Unit]
Description=http server v${${PROJECT_NAME}_VERSION_MAJOR}.${${PROJECT_NAME}_VERSION_MINOR}.${${PROJECT_NAME}_VERSION_PATCH} (project ${PROJECT_NAME})
Wants=node-ipt.service
After=network.target node-ipt.service

[Service]
Type=simple
ExecStart=/opt/smf/http -C/opt/smf/http_v${${PROJECT_NAME}_VERSION_MAJOR}.${${PROJECT_NAME}_VERSION_MINOR}.cfg
Restart=always
RestartSec=12

[Install]
WantedBy=multi-user.target
